# target/linux/nuclei/vega/base-files/etc/profile.d/userconfig.sh
#Auto generated, please don't change it!
export PS1='\[\e[1;32;40m\]\u@dev:$PWD> \[\e[0m\]'
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/root
